# systemd/cognihub.service
[Unit]
Description=CogniHub FastAPI Server
After=network.target

[Service]
Type=simple
User=small_ed
WorkingDirectory=/home/small_ed/cognihub
ExecStart=/usr/bin/python3 -m uvicorn app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target